# infra/financiall.service
[Unit]
Description=financiALL - servidor de importacao de notas fiscais
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/financiall
Environment=FINANCIALL_DB_PATH=/home/pi/financiall/data/financiall.db
Environment=FINANCIALL_UPLOAD_DIR=/home/pi/financiall/data/uploads
ExecStart=/home/pi/financiall/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=5

# Hardening leve, compativel com acesso a arquivos em WorkingDirectory
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
